# voicechat.service
[Unit]
Description=voicechat — headless dictation daemon (mic -> whispervulkan -> paste)
After=graphical-session.target whispervulkan.service
Wants=whispervulkan.service

[Service]
Type=simple
ExecStart=%h/.local/bin/voicechat run
Restart=on-failure
RestartSec=2
# By default voicechat records from the system default source. To dictate from a specific
# mic, uncomment and set this to a source name from `pactl list sources short`.
#Environment=VOICECHAT_SOURCE=alsa_input.example-source

[Install]
WantedBy=default.target
